# data/probe-client.service.default
[Unit]
Description=Probe client Service
After=network.target

[Service]
Type=simple
Environment="RUST_LOG=info"
Restart=on-failure
RestartSec=10s
ExecStart=/opt/probe-client/target/release/probe-client
WorkingDirectory=/opt/probe-client
DynamicUser=yes
KillSignal=SIGINT
TimeoutStartSec=2

[Install]
WantedBy=multi-user.target